# cmd/demo.crontab
# hello
* * * * * echo "Hello World!"
# Print date
* * * * * echo "Hi! at:$(date +%FT%T%z)"
* * * * * echo "中文"
* * * * * pwd
# 模拟耗时脚本
* * * * * echo "开始执行" && sleep 30 && echo "结束执行"
* * * * * node -e "console.log('hello')" && sleep 15 && echo "结束执行"

# 0-4,20-30,*/5 * * * * ls -alh .
0-4,20-30,*/2 * * * * ls -alh .
